# docker/standalone/rootfs/etc/s6-overlay/s6-rc.d/redis/run
#!/bin/bash
exec 2>&1

exec redis-server --bind 127.0.0.1 --port 6379 --daemonize no
