# docker/services/rootfs/etc/s6-overlay/s6-rc.d/cron/run
#!/bin/bash
exec 2>&1

cd /app
exec /root/.bun/bin/bun run applications/cron/dist/index.js
