# packages/services/rootfs/etc/s6-overlay/s6-rc.d/api/run
#!/bin/bash
exec 2>&1

cd /app
exec /root/.bun/bin/bun run packages/api/dist/index.js
