# docker/services/rootfs/etc/s6-overlay/s6-rc.d/web/run
#!/bin/bash
exec 2>&1

cd /app/applications/web
export HOSTNAME=0.0.0.0
export PORT=3000
exec /root/.bun/bin/bun dist/server-entry/index.js
